# puppet/modules/drupal/files/etc/cron.d/puppet
@reboot root /usr/site/lab/sbin/run_puppet 2>/dev/null
30 * * * * root /usr/site/lab/sbin/run_puppet 2>/dev/null
